# zerotier-one/rootfs/etc/services.d/zerotier-one/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# TomCzhen's Hass.io Add-ons: ZeroTier One
# Take down the S6 supervision tree when ZeroTier One fails
# ==============================================================================

source /usr/lib/hassio-addons/base.sh

api_port=$(hass.config.get 'api_port')

exec zerotier-one -U -p${api_port}